# lustre/changes/debian/lustre-#VERSION#.postrm
#!/bin/sh

if [ -x /usr/sbin/invoke-rc.d ] ; then
   update-rc.d -f #VERSION# remove || true
fi		
